# debian/php5-fpm.postinst
#!/bin/sh

set -e

if [ "$1" = "triggered" ] && [ "$2" = "/etc/php5/conf.d" ]; then
    invoke-rc.d php5-fpm restart
fi

dpkg-maintscript-helper mv_conffile /etc/php5/fpm/main.conf /etc/php5/fpm/php-fpm.conf 5.3.5-1 -- "$@"

if [ "$1" = "configure" ]; then
    phpini="/etc/php5/fpm/php.ini"

    ucf /usr/share/php5/php.ini-production $phpini
    ucfr php5-fpm $phpini
fi

#DEBHELPER#

exit 0
